# supervisor.conf
[program:scholar]
command=python /root/tornado_proxy/tornado_proxy/proxy.py
autorestart=True
